# etc/rc.d/iclcd
#!/bin/ksh
#
# $OpenBSD$

daemon="/usr/sbin/iclcd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
